# deploy/polytape-refresh.timer
# Runs the match-set refresh every ~10 minutes (and shortly after boot), so
# resolved matches roll out and newly-created fixtures roll in without manual
# intervention. Install alongside polytape-refresh.service:
#   sudo cp deploy/polytape-refresh.{service,timer} /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now polytape-refresh.timer
[Unit]
Description=Run polytape match-set refresh every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
