# deploy/polytape-refresh.service
# Oneshot: re-discover the open World Cup match set and restart polytape if it
# changed (rolls resolved matches out, new fixtures in). Driven by the timer.
# Runs as root (needs to write /etc/polytape and restart the service).
[Unit]
Description=polytape match-set refresh (roll out resolved, roll in new)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/polytape/polytape-refresh.sh
